# deploy/visual-buct.service
[Unit]
Description=Visual BUCT Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=yanp
Group=yanp
WorkingDirectory=/opt/visual_buct/backend
Environment="PATH=/opt/visual_buct/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/visual_buct/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4

# 自动重启
Restart=always
RestartSec=5

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=visual-buct

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
