# deploy/visual-buct.service
[Unit]
Description=Visual BUCT Backend API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/visual-buct/backend
Environment="PATH=/var/www/visual-buct/venv/bin"
ExecStart=/var/www/visual-buct/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2

# 自动重启
Restart=always
RestartSec=5

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=visual-buct

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
